# aconfmgr/.config/aconfmgr/old2/files/etc/profile.d/multi-gpu.sh
# multi gpu setup
export KWIN_DRM_DEVICES=/dev/dri/card1:/dev/dri/card0
export DRI_PRIME=pci-0000_03_00_0
