# aconfmgr/.config/aconfmgr/old2/files/etc/systemd/system/set-zcasdns.timer
[Unit]
Description=Set ZimmerCAS DNS on ZeroTier
Requires=set-zcasdns.service

[Timer]
Unit=set-zcasdns.service
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
