# aconfmgr/.config/aconfmgr/legacy/files/etc/systemd/system/sponsorblockcast.service.d/override.conf
[Service]
#Environment="SBCPOLLINTERVAL=10"
#Environment="SBCSCANINTERVAL=100"
Environment="SBCCATEGORIES=sponsor selfpromo interaction intro outro preview music_offtopic"